# 강의자료/실습및과제/tomcat/tomcat multi-instance/tomcat-hmo

#!/bin/sh 
# 
# Startup script for Tomcat for HMO
# 
# chkconfig: 35 85 35 
# description: Start Tomcat 
# 
# processname: tomcat 
# 
# Source function library. 

. /etc/rc.d/init.d/functions 

export JAVA_HOME=/usr/local/sunnyvale/jdk1.7
export CLASSPATH=.:$JAVA_HOME/lib/tools.jar
export CATALINA_HOME=/usr/local/sunnyvale/tomcat7
export PATH=$PATH:$JAVA_HOME/bin

case "$1" in 

	start) 

		echo -n "Starting tomcat-hmo: " 
		daemon /usr/local/sunnyvale/hmo/bin/startup.sh 
		touch /var/lock/subsys/tomcat-hmo 
		echo
		;; 
	stop) 
		echo -n "Shutting down tomca-hmo: " 
		daemon /usr/local/sunnyvale/hmo/bin/shutdown.sh 
		rm -f /var/lock/subsys/tomcat-hmo 
		echo 
		;; 
	restart) 
		$0 stop
		sleep 2 
		$0 start 
		;; 
	*) 
		echo "Usage: $0 {start|stop|restart}" 
		exit 1 
esac 
exit 0
